# public/service/google-calendar-backup/files/etc/cron.d/google-calendar-backup
%header("##")
MAILTO=root
0 1 * * * %(user) /usr/local/bin/google-calendar-backup && gzip -f9 %(data_path)/calendars/$(date +\%%Y-\%%m-\%%d)/*.ics
