# public/service/amavis/files/etc/cron.d/tmpreaper-for-virusmails
%header("##")
20 4    * * *   root    /usr/sbin/tmpreaper --mtime 30d %(quarantinedir)
